# rc.local
#!/bin/bash
#copie esse arquivo para /etc e garante que ele tenha direito de execução
#This service went to systemd
#/home/pdsilva/bin/fan.sh
/home/pdsilva/bin/pisca.sh

exit 0

